# clients/kindle/bin/slowdash.conf
start on started lab126_gui
stop on stopping lab126_gui
respawn
respawn limit 5 60

env SRCDIR=/mnt/us/extensions/slowdash/bin
env RUNDIR=/tmp/slowdash

pre-start script
    source /etc/upstart/functions

    # Copy scripts to tmpfs so the scheduler never holds open file
    # handles on /mnt/us.  The fsp FUSE daemon cannot cleanly unmount
    # the FAT32 partition for USB export while files are open, which
    # corrupts the FAT cluster chain.
    mkdir -p "${RUNDIR}"
    cp "${SRCDIR}"/*.sh "${RUNDIR}/"

    # Also copy fbink binary to /tmp if bundled with the extension
    if [ -x "${SRCDIR}/../fbink" ]; then
        cp "${SRCDIR}/../fbink" /tmp/fbink
        chmod +x /tmp/fbink
    fi
end script

exec /bin/sh /tmp/slowdash/scheduler.sh

post-stop script
    rm -rf /tmp/slowdash
end script
